# plugin/rootfs/lib/systemd/system/myexampledriver.socket
[Unit]
Description=myexampledriver

[Socket]
ListenStream=/run/docker/plugins/myexampledriver.sock

[Install]
WantedBy=sockets.target

